# webalizer/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command = /usr/sbin/apachectl -DFOREGROUND
autostart = true
autorestart = true
stopsignal=QUIT

